# pi-midi-http_0.1-1_armhf/DEBIAN/prerm
#!/bin/bash
sudo systemctl stop midi_http
sudo systemctl disable midi_http
sudo systemctl daemon-reload
